# postfixadmin.ini
[program:php-fpm]
command=php-fpm83

[program:nginx]
command=nginx